# playbook/templates/httpecho.service.j2
[Unit]
Description=httpecho

[Service]
ExecStart=/usr/bin/httpecho

[Install]
WantedBy=multi-user.target
